# docs/buyify.service
[Unit]
Description=Buyify

[Service]
ExecStart=/usr/bin/java -jar /home/buyify/buyify-0.0.1-SNAPSHOT.jar
User=buyify

[Install]
WantedBy=multi-user.target
